# scripts/admin/com.editorial-scripts.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.editorial-scripts.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dylanpossamai/Library/CloudStorage/Dropbox/Work/editorial_scripts/scripts/weekly_run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dylanpossamai/Library/CloudStorage/Dropbox/Work/editorial_scripts/results/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dylanpossamai/Library/CloudStorage/Dropbox/Work/editorial_scripts/results/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/dylanpossamai/Library/CloudStorage/Dropbox/Work/editorial_scripts</string>
</dict>
</plist>
